# docmost/docmost-pg-backup.timer
[Unit]
Description=Docmost Backup Timer

[Timer]
OnCalendar=*-*-* 2:00:00

[Install]
WantedBy=timers.target
