# home-pascal-navidrome-music.mount
# Automatically generated by systemd-fstab-generator

[Unit]
Documentation=man:fstab(5) man:systemd-fstab-generator(8)
SourcePath=/etc/fstab
Requires=network-online.target
After=network-online.target

[Mount]
What=//lil-nas.local/media/music
Where=/home/pascal/navidrome/music
Type=cifs
Options=credentials=/home/pascal/.smbcredentials,uid=www-data,gid=www-data,dir_mode=0775,file_mode=0664,soft,nofail,x-systemd.automount,x-systemd.requires=network-online.target,_netdev
