# immich/immich-immich-server.container
[Unit]
Requires=immich-redis.service immich-database.service
After=immich-redis.service immich-database.service

[Container]
ContainerName=immich_server
Image=ghcr.io/immich-app/immich-server:v2.4.0
Pod=immich.pod
Volume=/home/pascal/immich/library:/data
Volume=/etc/localtime:/etc/localtime:ro
Environment=DB_HOSTNAME=immich_postgres

[Service]
ExecStartPre=/bin/sleep 2
Restart=always
